# docker/cron/crontab
#!/bin/bash

*/5 * * * * root cd /var/www/app/ ; /usr/local/bin/php  bin/console app:update-newsfeed  >> /var/log/cron.log 2>&1




